# es-reindexers/systemd/perfscale-es-ols-load-test-results-ocp-qe.timer
[Unit]
Description=Perfscale ES OLS load test results ocp qe Timer
Requires=perfscale-es-ols-load-test-results-ocp-qe.service

[Timer]
OnCalendar=*-*-* *:0/8:00
Unit=perfscale-es-ols-load-test-results-ocp-qe.service

[Install]
WantedBy=timers.target